# versions/1.0-cron/config/cron/crontab
# /var/log/cron/cron.log
*/1 * * * * php composer --version